# supervisord.conf
[supervisord]
nodaemon=true

[program:mongod]
command=mongod --dbpath /data/db
redirect_stderr=true

[program:grunt]
command=grunt
redirect_stderr=true
